# packaging/debian/DEBIAN/postinst
#!/bin/sh
set -e

APP_ROOT=/opt/orbitdeck
VENV="$APP_ROOT/.venv"
PYTHON_BIN="${PYTHON_BIN:-python3}"

if [ ! -d "$VENV" ]; then
  "$PYTHON_BIN" -m venv "$VENV"
fi

"$VENV/bin/python" -m pip install --upgrade pip
"$VENV/bin/python" -m pip install --no-index --find-links "$APP_ROOT/vendor/wheels" -r "$APP_ROOT/requirements.txt"

if command -v systemctl >/dev/null 2>&1; then
  systemctl daemon-reload || true
  systemctl enable orbitdeck.service || true
  systemctl restart orbitdeck.service || systemctl start orbitdeck.service || true
fi
