# packaging/debian/usr/lib/systemd/system/orbitdeck.service
[Unit]
Description=OrbitDeck API
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/orbitdeck
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/orbitdeck/.venv/bin/python -m uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
